# template/supervisor_conf.d/app.ini
[program:app] ;
command = python /app/chooseClass.py ;
startsecs = 10        ;
redirect_stderr = true  ;
autorestart = true   ;
stdout_logfile_maxbytes = 50MB  ;
stdout_logfile = /app/app.log